# modules/init_files/templates/kube-controller-manager.service.tpl
[Unit]
Description=Kubernetes Controller Manager
Documentation=https://github.com/kubernetes/kubernetes

[Service]
ExecStart=/usr/local/bin/kube-controller-manager \
  --bind-address=0.0.0.0 \
  --cloud-provider="aws" \
  --cluster-cidr=${CLUSTER_CIDR} \
  --cluster-name=${CLUSTER_NAME} \
  --cluster-signing-cert-file=/var/lib/kubernetes/ca-cert.pem \
  --cluster-signing-key-file=/var/lib/kubernetes/ca-key.pem \
  --kubeconfig=/var/lib/kubernetes/kube-controller-manager.kubeconfig \
  --leader-elect=true \
  --root-ca-file=/var/lib/kubernetes/ca-cert.pem \
  --service-account-private-key-file=/var/lib/kubernetes/service-account-key.pem \
  --service-cluster-ip-range=${KUBERNETES_SERVICE_CIDR} \
  --use-service-account-credentials=true \
  --v=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target